# zfs-backup.cron.hourly
#!/bin/sh
exec /bin/bash zfs-backup -f /etc/default/zfs-backup